# infra/hermes-bridge/hermes-bridge.service
[Unit]
Description=T3 Code Hermes bridge
After=network-online.target hermes-webui.service
Wants=network-online.target

[Service]
Type=simple
User=forbes
EnvironmentFile=/home/forbes/.config/hermes-bridge.env
ExecStart=/usr/bin/python3 /home/forbes/hermes-bridge/hermes_bridge.py
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
